# .ebextensions/fstrim.cron
*/30 * * * * sudo fstrim /proc/$(sudo docker inspect --format='{{ .State.Pid }}' $(sudo docker ps -q))/root
